# spark-monitoring.service
[Unit]
Description=Stream vLLM container logs and GPU stats to host filesystem
After=docker.service
Requires=docker.service

[Service]
Type=simple
User=__USER__
Restart=always
RestartSec=10
ExecStart=/home/__USER__/spark_monitoring_tui/spark-monitoring.sh

[Install]
WantedBy=multi-user.target
